# pi-appliance/systemd/tp-signer-kiosk.service
[Unit]
Description=TP Offline Signer Kiosk
After=multi-user.target pcscd.service
Wants=pcscd.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/tp-pi-kiosk
Environment=TP_KIOSK_WORK_DIR=/run/tp-pi-kiosk
RuntimeDirectory=tp-pi-kiosk
ExecStart=/opt/tp-pi-kiosk/launch_kiosk.sh
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
